# scripts/com.commonplace.audiobooks-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.commonplace.audiobooks-scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/cameronlewis/code/commonplace-code/.venv/bin/python</string>
        <string>/Users/cameronlewis/code/commonplace-code/scripts/audiobooks_scan.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/cameronlewis/code/commonplace-code</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>COMMONPLACE_DB_PATH</key>
        <string>/Users/cameronlewis/commonplace/library.db</string>
    </dict>

    <!-- Run daily at 04:00 — catches new audiobooks overnight.
         Scanner is idempotent (skips already-ingested) and exits
         cleanly with logged error when the external drive is unmounted. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/cameronlewis/Library/Logs/commonplace-audiobooks-scan.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/cameronlewis/Library/Logs/commonplace-audiobooks-scan.err.log</string>
</dict>
</plist>
